# template/cronweb.service.tmpl
[Unit]
Description=CronWeb daemon

[Service]
Group={group}
User={user}
WorkingDirectory={pwd}
Restart=always
Type=simple
ExecStart={exec}

[Install]
WantedBy=multi-user.target
